# deploy/systemd/img-studio-xai-auth-sync.timer
[Unit]
Description=Refresh IMG Studio xAI OAuth pool every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=img-studio-xai-auth-sync.service

[Install]
WantedBy=timers.target
